# forex-bot.service
[Unit]
Description=Forex Demo Trading Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace
Environment=PATH=/workspace/venv/bin
Environment=PYTHONPATH=/workspace
Environment=PYTHONUNBUFFERED=1
Environment=TELEGRAM_BOT_TOKEN=YOUR_TELEGRAM_BOT_TOKEN
Environment=TELEGRAM_CHAT_ID=YOUR_TELEGRAM_CHAT_ID
Environment=OANDA_API_KEY=YOUR_OANDA_API_KEY
Environment=OANDA_ACCOUNT_ID=YOUR_OANDA_ACCOUNT_ID
ExecStart=/workspace/venv/bin/python bot_runner.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=forex-bot

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

LimitNOFILE=65536
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target